# contrib/fetch-all-tools.timer
[Unit]
Description=AI Token Analyzer - Auto-detect and Fetch Tool Data Timer
Requires=fetch-all-tools.service

[Timer]
OnBootSec=60s
OnUnitActiveSec=30s
AccuracySec=1s
Persistent=true

[Install]
WantedBy=timers.target
