# contrib/fetch-all-tools.service
[Unit]
Description=AI Token Analyzer - Auto-detect and Fetch Tool Data
After=network.target

[Service]
Type=simple
User=openclaw
Group=openclaw
WorkingDirectory=/opt/ai-token-analyzer
Environment="PYTHONUNBUFFERED=1"
Environment="UPLOAD_HOSTNAME=ai-lab"
Environment="FETCH_DAYS=1"
ExecStart=/usr/bin/python3 /opt/ai-token-analyzer/scripts/fetch_all_tools.py

[Install]
WantedBy=multi-user.target
