# contrib/upload-service.service
[Unit]
Description=AI Token Analyzer - Upload Service
After=network.target

[Service]
Type=simple
Restart=always
RestartSec=5
User=openclaw
Group=openclaw
WorkingDirectory=/opt/ai-token-analyzer
Environment="PYTHONUNBUFFERED=1"
Environment="UPLOAD_SERVER=http://192.168.31.181:5001"
Environment="UPLOAD_AUTH_KEY=deploy-remote-machine-key-2026"
Environment="UPLOAD_HOSTNAME=ai-lab"
ExecStart=/usr/bin/python3 /opt/ai-token-analyzer/scripts/upload_to_server.py \
    --server $UPLOAD_SERVER \
    --auth-key $UPLOAD_AUTH_KEY \
    --hostname $UPLOAD_HOSTNAME \
    --daemon \
    --interval 30 \
    --incremental
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ai-token-upload

[Install]
WantedBy=multi-user.target
